# etc/systemd/system/router-healthcheck.service
[Unit]
Description=Router auto-recovery watchdog

[Service]
ExecStart=/usr/local/bin/router-healthcheck.sh
